# i3/PKG/post-install
#!/bin/sh

systemctl --user enable greenclip.service
systemctl --user start greenclip.service
